# systemd/gitup.service
[Unit]
Description=GitUp git auto commits

[Service]
Type=oneshot
ExecStart=/opt/minecraft-server/server/tekxit2-server/gitup.sh
